# rootfs/etc/s6-overlay/s6-rc.d/discourse-bootstrap/up
/etc/entrypoint.d/20-discourse-bootstrap.sh
